# crontab
0 8 * * 1 cd /app && python manage.py generate_weekly_content 2>&1
*/30 * * * * cd /app && python manage.py publish_due_posts 2>&1
