# crontab
# Monday 8:00 — generate text posts for next week (fast, ~30s)
0 8 * * 1 cd /app && python manage.py generate_weekly_content 2>&1

# Monday 8:05 — generate images in parallel for any posts missing them (4 workers, ~5 min total)
5 8 * * 1 cd /app && python manage.py generate_missing_images 2>&1

# Monday 8:35 — generate videos in parallel for any posts missing them (3 workers, ~8 min total)
35 8 * * 1 cd /app && python manage.py generate_missing_videos --all 2>&1

# Every 30 min — publish approved posts that are due
*/30 * * * * cd /app && python manage.py publish_due_posts 2>&1
